# esp_sync/esp_sync.timer
[Unit]
Description=Run esp_sync

[Timer]
OnStartupSec=10

[Install]
WantedBy=timers.target
