# dnsmasq_restarter/dnsmasq_restarter.path
[Path]
Unit=dnsmasq_restarter.service
PathModified=/etc/dnsmasq.d

[Install]
WantedBy=multi-user.target
